# config/frpc.supervisor.conf
[program:center-frpc-ssh]
command=/root/center_frpc_ssh/bin/frpc
restart=always
user=root
