# config/frpc.supervisor.conf
[program:frpc]
command=/root/center_frpc_ssh/bin/frpc
restart=always
user=root
