# config/frpc.supervisor.conf
[program:center-frpc]
command=/root/center_frpc/bin/frpc
stopasgroup=true
killasgroup=true
restart=always
user=root
